# docker/s6-rc.d/ssr/run
#!/command/with-contenv sh

cd /var/www

php artisan inertia:start-ssr
