# docker/s6-rc.d/symfony/up
/etc/s6-overlay/s6-rc.d/symfony/init.sh
